# tools/scripts/ir-daemon.service
[Unit]
Description=IR-Trigger Daemon for AD00020P
After=network.target

[Service]
Type=simple
User=root
# Replace /PATH/TO/... with your actual directory paths
# Example: /home/pi/IR-Trigger/edge_scripts/venv/bin/python
ExecStart=/PATH/TO/PYTHON /PATH/TO/ir_daemon.py --url http://YOUR_HA_IP:8123/api/webhook/ir_trigger_webhook --receiver MyReceiver
Restart=always
RestartSec=5
# If pyusb requires root privileges to access the USB device, keep User=root
# Alternatively, configure udev rules and run as a non-root user.

[Install]
WantedBy=multi-user.target
