# edge_scripts/ir-daemon.service
[Unit]
Description=IR-Trigger Daemon for AD00020P
After=network.target

[Service]
Type=simple
User=root
# Replace with the actual path to the script and your Home Assistant URL
ExecStart=/usr/bin/python3 /opt/ir_trigger/ir_daemon.py --url http://127.0.0.1:8123/api/webhook/ir_trigger_webhook --receiver LivingRoom_AD00020P
Restart=always
RestartSec=5
# If pyusb requires root privileges to access the USB device, keep User=root
# Alternatively, configure udev rules and run as a non-root user.

[Install]
WantedBy=multi-user.target
